# crontab
* * * * * sh -c "java -cp /app/target/application.jar ar.edu.utn.frba.dds.impactoambiental.JobEnvioGuia"
